# templates/default/plugins/nginx-vhosts/nginx-reloader.conf
description     "Dokku nginx reloader service"

start on runlevel [2345]
stop on runlevel [!2345]

script
  [[ -f /home/dokku/reload-nginx ]] && rm -rf /home/dokku/reload-nginx
  echo | sudo -u dokku nc -l -U /home/dokku/reload-nginx && /etc/init.d/nginx reload
end script
respawn
