# templates/default/plugins/nginx-vhosts/nginx-reloader.conf
description     "Dokku nginx reloader service"

start on runlevel [2345]
stop on runlevel [!2345]

script
  [[ -f /home/git/reload-nginx ]] && rm -rf /home/git/reload-nginx
  echo | sudo -u git nc -l -U /home/git/reload-nginx && /etc/init.d/nginx reload
end script
respawn
